# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── Redis ────────────────────────────────────────────────────
[program:redis]
command=/usr/bin/redis-server --bind 127.0.0.1 --port 6379 --save "" --appendonly no --loglevel warning
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/redis.log
stderr_logfile=/var/log/redis.log
stdout_logfile_maxbytes=5MB

; ── PostgreSQL ───────────────────────────────────────────────
[program:postgres]
command=/usr/lib/postgresql16/bin/postgres -D /var/lib/postgresql/data -k /run/postgresql
user=postgres
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/postgres.log
stderr_logfile=/var/log/postgres.log
stdout_logfile_maxbytes=5MB

; ── Go API ───────────────────────────────────────────────────
[program:api]
command=/app/api/api-server
directory=/app/api
autostart=true
autorestart=true
priority=40
startsecs=5
startretries=10
stdout_logfile=/var/log/api.log
stderr_logfile=/var/log/api.log
stdout_logfile_maxbytes=10MB
environment=PORT="%(ENV_PORT)s",DB_HOST="%(ENV_DB_HOST)s",DB_PORT="%(ENV_DB_PORT)s",DB_USER="%(ENV_DB_USER)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s",DB_NAME="%(ENV_DB_NAME)s",DB_SSLMODE="%(ENV_DB_SSLMODE)s",JWT_SECRET="%(ENV_JWT_SECRET)s",REDIS_ADDR="%(ENV_REDIS_ADDR)s"

; ── Next.js Web ──────────────────────────────────────────────
[program:web]
command=/usr/bin/node server.js
directory=/app/web
autostart=true
autorestart=true
priority=50
startsecs=5
startretries=5
stdout_logfile=/var/log/web.log
stderr_logfile=/var/log/web.log
stdout_logfile_maxbytes=10MB
environment=NODE_ENV="production",PORT="5400",HOSTNAME="0.0.0.0",NEXT_PUBLIC_API_URL="%(ENV_NEXT_PUBLIC_API_URL)s"
